# data/cobbled.service
[Unit]
Description=Cobble — Pebble BLE daemon
Documentation=https://github.com/luigi311/libpebble-ble
After=bluetooth.target graphical-session.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/cobbled
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
