# data/pebble-led.service
[Unit]
Description=Pebble BLE daemon
Documentation=https://github.com/luigi311/libpebble-ble
After=bluetooth.target graphical-session.target
Wants=bluetooth.target
# Skip the service cleanly when the config file hasn't been created yet; without
# this, an absent file causes ExecStart to fail, triggering a Restart=on-failure loop.
ConditionPathExists=%E/pebble-led/config.toml

[Service]
Type=simple
ExecStart=/usr/bin/pebble-led
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
