# packages/pebble-led/data/pebble-led.service
[Unit]
Description=Pebble BLE daemon (libpebble-ble)
After=bluetooth.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
EnvironmentFile=%h/.config/pebble-led/env
ExecStart=/usr/bin/pebble-led ${PEBBLE_ADDRESS}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
