# packages/pebble-led/data/pebble-led.service
[Unit]
Description=Pebble LE daemon (owns the BLE connection to the watch)
Documentation=https://github.com/yourname/pebble-le
After=bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
# Set your watch address here, or override with `systemctl --user edit`.
Environment=PEBBLE_WATCH_ADDRESS=E6:94:0A:D4:D5:DC
ExecStart=pebble-led ${PEBBLE_WATCH_ADDRESS}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
